# discoveryengine/supervisord.conf
[supervisord]

[program:uwsgi]
command=/home/ec2-user/.virtualenvs/discovery/bin/uwsgi --ini /home/ec2-user/discovery-engine/discoveryengine/uwsgi.ini
user=ec2-user